# config/docker/web/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:apache2]
command=/bin/bash -c "/var/www/config/docker/web/bootstrap.sh"

[program:mailcatcher]
command=/bin/bash -c "if [[ $ENVIRONMENT == 'DEV' ]]; then /usr/bin/env $(which mailcatcher) --foreground --http-ip=0.0.0.0; fi"
autorestart = false

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true
